# taobei/challenge-06/supervisord.conf
[unix_http_server]
file=/tmp/supervisor-taobei.sock

[supervisord]
logfile=/tmp/supervisord-taobei.log
pidfile=/tmp/supervisord-taobei.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-taobei.sock

[program:tbbuy]
command=python -u -m tbbuy.app
stdout_logfile=/tmp/tbbuy.log
stderr_logfile=/tmp/tbbuy.log

[program:tbfile]
command=python -u -m tbfile.app
stdout_logfile=/tmp/tbfile.log
stderr_logfile=/tmp/tbfile.log

[program:tbmall]
command=python -u -m tbmall.app
stdout_logfile=/tmp/tbmall.log
stderr_logfile=/tmp/tbmall.log

[program:tbuser]
command=python -u -m tbuser.app
stdout_logfile=/tmp/tbuser.log
stderr_logfile=/tmp/tbuser.log

[program:tbweb]
command=python -u -m tbweb.app
stdout_logfile=/tmp/tbweb.log
stderr_logfile=/tmp/tbweb.log
